# scripts/picrawler-mcp.service
[Unit]
Description=PiCrawler MCP server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=smaniak
WorkingDirectory=/home/smaniak/picrawler-app
ExecStart=/home/smaniak/picrawler-app/.venv/bin/python /home/smaniak/picrawler-app/mcp_server.py
Restart=on-failure
RestartSec=3
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
